# rpm/daq_token.socket
[Unit]
Description=The JWT token meister socket.

[Socket]
ListenStream=%t/daq_token
