# etc/systemd/system/tdaq_token_gssapi.socket
[Unit]
Description=The JWT token meister socket using gssapi.

[Socket]
ListenStream=8990
