# etc/systemd/system/tdaq_token.socket
[Unit]
Description=The JWT token meister socket.

[Socket]
SocketUser=tdaqsw
ListenStream=%t/tdaq_token
